# transcode-node/supervisord.conf
[supervisord]
nodaemon=true
logfile=/app/logs/supervisord.log
pidfile=/var/run/supervisord.pid

[program:transcode-worker]
command=python3 -m app.main
directory=/app
autostart=true
autorestart=true
stderr_logfile=/app/logs/worker.err.log
stdout_logfile=/app/logs/worker.out.log
